# scripts/systemd/startup.service
[Unit]
Description=Startup Script

[Service]
Type=oneshot
ExecStart=/etc/scripts/startup

[Install]
WantedBy=multi-user.target
